# dockers/salt-master/supervisor-salt.conf
[program:salt-master]
command=/usr/bin/salt-master
numprocs=1
autostart=true
autorestart=true

[program:salt-api]
command=/usr/bin/salt-api
numprocs=1
autostart=true
autorestart=true

