# dockers/salt-master/supervisord.conf
[supervisord]
nodaemon=true
minprocs=50

[program:salt-master]
command=/usr/bin/salt-master
autorestart=true
priority=10

[program:salt-api]
command=/usr/bin/salt-api
autorestart=true
priority=20
