# dockers/salt-minion/supervisord.conf
[supervisord]
nodaemon=true

[program:salt-minion]
command=/usr/bin/python /usr/bin/salt-minion
autorestart=true
priority=10